# system/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#
hdparm -B 254 /dev/sda >> /dev/null
